# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
loglevel=info
pidfile=/var/run/supervisord.pid

[program:openclaw-api]
command=python3 -m uvicorn core.main:app --host 0.0.0.0 --port 8080 --workers 1
directory=/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/log/openclaw-api.log
stderr_logfile=/var/log/openclaw-api.log
environment=PYTHONPATH="/app"

[program:discord-bot]
command=python3 /app/bots/discord_bot.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/discord.log
stderr_logfile=/var/log/discord.log
environment=PYTHONPATH="/app"

[program:slack-bot]
command=python3 /app/bots/slack_bot.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/slack.log
stderr_logfile=/var/log/slack.log
environment=PYTHONPATH="/app"

[program:telegram-bot]
command=python3 /app/bots/telegram_bot.py
directory=/app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/telegram.log
stderr_logfile=/var/log/telegram.log
environment=PYTHONPATH="/app"

[program:whatsapp-bot]
command=node /app/bots/whatsapp/index.js
directory=/app/bots/whatsapp
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/whatsapp.log
stderr_logfile=/var/log/whatsapp.log

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface
